# debian/jellyfin.upstart
description "jellyfin daemon"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

console log
respawn
respawn limit 10 5

kill timeout 20

script
	set -x
	echo "Starting $UPSTART_JOB"

	# Log file
	logger -t "$0" "DEBUG: `set`"
    . /etc/default/jellyfin
    exec su -u $JELLYFIN_USER -c /usr/bin/jellyfin $JELLYFIN_ARGS
end script
